# zerokube/zerokube.service
[Unit]
Description=Zerokube
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zerokube
Restart=on-failure
